# deploy/hospital-monitor.service
[Unit]
Description=Sistema de Monitoramento de Quedas Hospitalares
# Aguarda a sessão gráfica estar pronta (importante para cv2.imshow)
After=graphical-session.target network.target
Wants=graphical-session.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/hospital-monitor
ExecStart=/usr/bin/python3 /opt/hospital-monitor/main.py
Restart=always
RestartSec=5
StandardOutput=append:/var/log/hospital-monitor/service.log
StandardError=append:/var/log/hospital-monitor/service.log

# Watchdog do systemd (service deve enviar heartbeat)
WatchdogSec=120

# Variaveis de ambiente para GUI
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority

# Timeout maior para aguardar GUI inicializar
TimeoutStartSec=180

# Limites de recursos
MemoryMax=512M
CPUQuota=80%

[Install]
# Usar graphical-session.target para garantir que GUI está pronta
WantedBy=graphical-session.target
