# deploy/hospital-monitor-headless.service
[Unit]
Description=Sistema de Monitoramento de Quedas Hospitalares MVision (Headless)
After=network.target
Wants=network.target

# Remove limite de reinicializações - SEMPRE reinicia
StartLimitIntervalSec=0

[Service]
Type=simple
User=tmed
Group=tmed

# Diretório do código
WorkingDirectory=/mvision
ExecStart=/usr/bin/python3 -u /mvision/main.py

# SEMPRE reinicia, não importa o motivo da falha
Restart=always
RestartSec=10

# Logs
StandardOutput=journal
StandardError=journal

# Variáveis de ambiente (sem DISPLAY para modo headless)
Environment=PYTHONUNBUFFERED=1
Environment=QT_QPA_PLATFORM=offscreen

# Timeout generoso para inicialização (YOLO demora para carregar)
TimeoutStartSec=300
TimeoutStopSec=30

# Limites de recursos (YOLO precisa de memória)
MemoryMax=768M
CPUQuota=90%

# Permite acesso a GPIO e câmera
SupplementaryGroups=gpio video

[Install]
WantedBy=multi-user.target
